# cookbooks/bcpc/templates/default/ceph-upstart.ceph-mon-renice.conf.erb
description "Ceph monitor renicer"

start on started ceph-mon-all

post-start script
  pgrep ceph-mon | xargs renice <%=node['bcpc']['ceph']['mon_niceness']%>
end script

